# pika-baseos/debian/pika-baseos.postinst
#! /bin/bash
systemctl disable --now gpu-manager
systemctl mask gpu-manager
systemctl enable --now nvidia-pika-prime