# pika-baseos/debian/pika-baseos.prerm
#! /bin/bash
systemctl disable --now nvidia-pika-prime
systemctl unmask gpu-manager
systemctl enable --now gpu-manager
