# pika-baseos/debian/pika-baseos.postinst
#! /bin/bash
systemctl disable --now gpu-manager || true
systemctl mask gpu-manager || true
systemctl enable --now nvidia-pika-prime || true
