# pika-baseos/debian/pika-baseos.prerm
#! /bin/bash
systemctl disable --now nvidia-pika-prime || true
systemctl unmask gpu-manager || true
systemctl enable --now gpu-manager || true
