# pika-baseos/etc/profile.d/0-pika-path.sh
#! /bin/sh
export PATH="/usr/local/sbin:/usr/local/bin:/usr/local/libexec:/usr/sbin:/usr/bin:/sbin:/bin:/usr/libexec:/usr/local/games:/usr/games:$PATH"